# packaging/deb/systemd/code-remote-ngrok.service
[Unit]
Description=claude-remote ngrok edge (reserved domain -> Synapse :8080)
After=code-remote-synapse.service network-online.target
Wants=code-remote-synapse.service network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/code-remote/.env
ExecStart=/usr/bin/ngrok-forward
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
